# stunnel/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/stunnel
drwxr-xr-x	root/root	etc/ssl/
drwxr-xr-x	root/root	etc/ssl/certs/
-rw-r--r--	root/root	etc/ssl/certs/stunnel.crt (EMPTY)
drwxr-xr-x	root/root	etc/ssl/keys/
-rw-------	root/root	etc/ssl/keys/stunnel.key (EMPTY)
drwxr-xr-x	root/root	etc/stunnel/
-rw-r--r--	root/root	etc/stunnel/stunnel.conf
-rw-r--r--	root/root	etc/stunnel/stunnel.conf-sample
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/stunnel/
-rwxr-xr-x	root/root	usr/lib/stunnel/libstunnel.la
-rwxr-xr-x	root/root	usr/lib/stunnel/libstunnel.so
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/stunnel
-rwxr-xr-x	root/root	usr/sbin/stunnel3
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/stunnel.8.gz
-rw-r--r--	root/root	usr/share/man/man8/stunnel.pl.8.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxrwx--T	root/nobody	var/lib/stunnel/
